# deploy/systemd/defenseclaw-axis.target
[Unit]
Description=DefenseClaw mandatory AXIS services
Wants=defenseclaw-axis-broker.service
