# packaging/systemd/defenseclaw-agent-control.service
[Unit]
Description=DefenseClaw Agent Control Policy Synchronizer
Documentation=https://docs.defenseclaw.ai/docs/setup/agent-control
After=network-online.target defenseclaw-gateway.service
Wants=network-online.target defenseclaw-gateway.service

[Service]
Type=simple
User=defenseclaw
Group=defenseclaw
Environment=DEFENSECLAW_HOME=/var/lib/defenseclaw
Environment=DEFENSECLAW_CONFIG=/etc/defenseclaw/config.yaml
Environment=DEFENSECLAW_DEPLOYMENT_MODE=managed_enterprise
EnvironmentFile=-/etc/defenseclaw/agent-control.env
ExecStart=/opt/defenseclaw/bin/defenseclaw agent-control sync --watch
Restart=on-failure
RestartSec=5s
StateDirectory=defenseclaw
RuntimeDirectory=defenseclaw
LogsDirectory=defenseclaw
StateDirectoryMode=0750
RuntimeDirectoryMode=0750
LogsDirectoryMode=0750
WorkingDirectory=/var/lib/defenseclaw
UMask=0077
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectClock=true
ProtectControlGroups=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProcSubset=pid
ReadOnlyPaths=/etc/defenseclaw /opt/defenseclaw
ReadWritePaths=/var/lib/defenseclaw /var/log/defenseclaw /run/defenseclaw
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
