# packaging/launchd/com.defenseclaw.hook-guardian.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.defenseclaw.hook-guardian</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Library/DefenseClaw/bin/defenseclaw-gateway</string>
    <string>enterprise</string>
    <string>hooks</string>
    <string>reconcile</string>
    <string>--manifest</string>
    <string>/Library/Application Support/DefenseClaw/hook-guardian/targets.yaml</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Library/Application Support/DefenseClaw</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>DEFENSECLAW_HOME</key>
    <string>/Library/Application Support/DefenseClaw</string>
    <key>DEFENSECLAW_CONFIG</key>
    <string>/Library/Application Support/DefenseClaw/config.yaml</string>
    <key>DEFENSECLAW_DEPLOYMENT_MODE</key>
    <string>managed_enterprise</string>
    <key>DEFENSECLAW_HOOK_GUARDIAN_AUTH_DIR</key>
    <string>/Library/Application Support/DefenseClaw/hook-guardian-state</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>Umask</key>
  <integer>63</integer>
  <key>StandardOutPath</key>
  <string>/Library/Logs/DefenseClaw/hook-guardian.log</string>
  <key>StandardErrorPath</key>
  <string>/Library/Logs/DefenseClaw/hook-guardian.err.log</string>
</dict>
</plist>
